# etc/system/systemd/parser.service
[Unit]
Description=SMTPHook Parser
After=network.target
Requires=smtphook.target

[Service]
Type=simple
User=smtphook
WorkingDirectory=/opt/smtphook/parser
EnvironmentFile=/opt/smtphook/parser/.env
ExecStart=/opt/smtphook/bin/parser
Restart=always
RestartSec=5s

[Install]
WantedBy=smtphook.target
