# etc/system/systemd/smtphook.target
[Unit]
Description=SMTPHook Service Stack
Requires=parser.service webhook-server.service
After=network.target

[Install]
WantedBy=multi-user.target
