# etc/system/systemd/webhook-server.service
[Unit]
Description=SMTPHook Webhook Receiver
After=network.target
Requires=smtphook.target

[Service]
Type=simple
User=smtphook
WorkingDirectory=/opt/smtphook/webhook-server
EnvironmentFile=/opt/smtphook/webhook-server/.env
ExecStart=/opt/smtphook/bin/webhook-server
Restart=always
RestartSec=5s

[Install]
WantedBy=smtphook.target
