# etc/system/systemd/webhook-server.service
[Unit]
Description=SMTPHook Webhook Receiver
After=network.target

[Service]
Type=simple
ExecStart=/opt/smtphook/bin/webhook-server
WorkingDirectory=/opt/smtphook/webhook-server
EnvironmentFile=/opt/smtphook/webhook-server/.env
Restart=always
User=nobody
Group=nogroup

[Install]
WantedBy=smtphook.target
